# conf/pixivdaily.timer
[Unit]
Description=Run PixivDaily bot every midnight

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true

[Install]
WantedBy=timers.target
